# deploy/igris.service
[Unit]
Description=IGRIS GPT Service
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=igris
WorkingDirectory=/home/igris/IGRIS_GPT
ExecStart=/home/igris/IGRIS_GPT/.venv/bin/python -m igris.web.server
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=PROJECT_ROOT=/home/igris/IGRIS_GPT
Environment=IGRIS_RANK_MAX_STEPS=80
Environment=IGRIS_LARGE_MISSION_REASONING_TIMEOUT=180
Environment=IGRIS_BASELINE_CACHE_SECONDS=3600
Environment=IGRIS_FORCE_PREEMPTIVE_DECOMPOSITION=true
Environment=IGRIS_PREFER_DETERMINISTIC_DECOMPOSITION=true

[Install]
WantedBy=multi-user.target
