# config/igris-gpt.service.example
[Unit]
Description=IGRIS_GPT AI Engineering Agent
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=igris
Group=igris
WorkingDirectory=/opt/igris-gpt
EnvironmentFile=/opt/igris-gpt/.env
Environment=IGRIS_RANK_MAX_STEPS=80
Environment=IGRIS_LARGE_MISSION_REASONING_TIMEOUT=180
Environment=IGRIS_BASELINE_CACHE_SECONDS=3600
Environment=IGRIS_FORCE_PREEMPTIVE_DECOMPOSITION=true
Environment=IGRIS_PREFER_DETERMINISTIC_DECOMPOSITION=true
ExecStart=/opt/igris-gpt/.venv/bin/python -m uvicorn igris.web.server:app --factory --host 0.0.0.0 --port 7778 --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/igris-gpt/logs/igris.log
StandardError=append:/opt/igris-gpt/logs/igris.log

[Install]
WantedBy=multi-user.target
